# templates/panel.crontab
# Sculptor scheduler do not modify
* * * * * sculptor schedule:run >> /dev/null 2>&1
