# templates/panel.crontab
0 5 * * 7 sudo certbot renew --nginx --non-interactive --post-hook "systemctl restart nginx.service"

