# templates/panel.crontab
# www-data scheduler do not modify
* * * * *  /usr/bin/php /var/www/html/current/artisan schedule:run >> /dev/null 2>&1
